# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=-/sbin/agetty --noreset --noclear --autologin holoiso - ${TERM}
